# docker/fernwaerts/s6/vector-log/run
#!/command/execlineb -P
s6-log -d3 T n10 s10000000 /var/log/vector
