# docker/fernwaerts/s6/auth/run
#!/command/with-contenv sh
exec 2>&1
export GOTRUE_API_HOST=127.0.0.1
export GOTRUE_API_PORT=9999
export GOTRUE_DB_DRIVER=postgres
export GOTRUE_DB_DATABASE_URL="postgres://supabase_auth_admin:${POSTGRES_PASSWORD}@${POSTGRES_HOST}:${POSTGRES_PORT}/${POSTGRES_DB}?sslmode=disable"
export GOTRUE_JWT_SECRET="${JWT_SECRET}"
export GOTRUE_JWT_EXP="${JWT_EXPIRY}"
export GOTRUE_JWT_ADMIN_ROLES=service_role
export GOTRUE_JWT_AUD=authenticated
export GOTRUE_JWT_DEFAULT_GROUP_NAME=authenticated
export GOTRUE_JWT_ISSUER="${API_EXTERNAL_URL}/auth/v1"
export GOTRUE_SITE_URL="${SITE_URL}"
export GOTRUE_DISABLE_SIGNUP="${DISABLE_SIGNUP}"
export GOTRUE_EXTERNAL_EMAIL_ENABLED="${ENABLE_EMAIL_SIGNUP}"
export GOTRUE_MAILER_AUTOCONFIRM="${ENABLE_EMAIL_AUTOCONFIRM}"
export GOTRUE_SMTP_ADMIN_EMAIL="${SMTP_ADMIN_EMAIL}"
export GOTRUE_SMTP_HOST="${SMTP_HOST}"
export GOTRUE_SMTP_PORT="${SMTP_PORT}"
export GOTRUE_SMTP_USER="${SMTP_USER}"
export GOTRUE_SMTP_PASS="${SMTP_PASS}"
export GOTRUE_SMTP_SENDER_NAME="${SMTP_SENDER_NAME}"
export API_EXTERNAL_URL="${API_EXTERNAL_URL}"
exec /opt/auth/auth
