# docker/fernwaerts/s6/meta-log/run
#!/command/execlineb -P
s6-log T n10 s10000000 /var/log/meta
