# supabase/docker/fernwaerts/s6/vector/run
#!/command/with-contenv sh
exec 2>&1
# vector.yml's sinks authenticate to Logflare with LOGFLARE_API_KEY; ingestion
# uses the public token.
export LOGFLARE_API_KEY="${LOGFLARE_PUBLIC_ACCESS_TOKEN}"
exec /usr/local/bin/vector --config /etc/vector/vector.yml
